# embedded-linux/K26.1/pham-dac-tam/session-06/Exercise_2/logger.service
[Unit]
Description=Logger Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/logger
StandardOutput=journal
StandardError=journal
SyslogIdentifier=logger

[Install]
WantedBy=multi-user.target
