# embedded-linux/K26.1/pham-dac-tam/session-06/Exercise_1/monitor.service
[Unit]
Description=Monitor Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/monitor
Restart=on-failure
RestartSec=5
StartLimitInterval=30
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
